# deps/go_ip_updater.service
[Unit]
Description=Route53 auto updater written in golang

[Service]
ExecStart=/usr/bin/rkt run \
	--volume conf,kind=host,source=/go_ip_updater/conf/goip.list \
  --mount volume=conf,target=/goip.list \
	unixvoid.com/go_ip_updater -- -loglevel debug

[Install]
WantedBy=multi-user.target
